# systemd.service.template
[Service]
ExecStart=<path to the executable>
WorkingDirectory=<dir>
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
